# services/kibana.service.conf
[program:kibana]
user=docker
command=/opt/kibana/bin/kibana
stdout_events_enabled = true
stderr_events_enabled = true
